# config/uwsgi.sff.conf
description "uwsgi processes"
start on runlevel [2345]
stop on runlevel [06]

exec uwsgi /var/django/sff/uwsgi.ini
